# ansible/roles/haproxy/templates/haproxy.service.j2
[Unit]
Description=HAProxy Load Balancer

[Service]
ExecStart=/usr/local/sbin/haproxy \
  -f /etc/haproxy/haproxy.cfg \
  -p /run/haproxy.pid $OPTIONS

[Install]
WantedBy=multi-user.target